# sidecar-oracle/systemd/awards-reconcile.service
[Unit]
Description=Awards Pipeline — Weekly reconciliation audit (sidecar)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=awards
Group=awards
WorkingDirectory=/opt/awards-pipeline/sidecar-oracle
EnvironmentFile=/opt/awards-pipeline/sidecar-oracle/.env
ExecStart=/usr/bin/node /opt/awards-pipeline/sidecar-oracle/run-reconcile.mjs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=awards-reconcile
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
MemoryMax=128M
TimeoutStartSec=10min
