# sidecar-oracle/systemd/awards-sam-vendors.timer
[Unit]
Description=Daily SAM.gov vendor enrichment
Requires=awards-sam-vendors.service

[Timer]
OnCalendar=*-*-* 10:00:00 UTC
Persistent=true
RandomizedDelaySec=15min
Unit=awards-sam-vendors.service

[Install]
WantedBy=timers.target
